# .docksal/services/cli/crontab
* * * * * curl -s -o /dev/null "http://uwec.docksal.site/scheduler/cron/14e040207abdf5e86fb1"